# systemd/hdparm.service
[Unit]
Description=hdparm sleep

[Service]
Type=oneshot
ExecStart=/sbin/hdparm -q -S 60 -y /dev/sda1

[Install]
WantedBy=multi-user.target
